# images/aquarium/root/etc/udev/rules.d/10-aqr-nics.rules
ACTION=="add", SUBSYSTEM=="net", RUN+="/usr/share/aquarium/boot/add-nic.sh"
